# crontab
MAIL=

29,59 * * * * sleep 45; python3 main.py
